# rootfs/rootfs/etc/systemd/system/wazuh-agent.service.d/override.conf
[Unit]
After=network.target
Wants=network.target

[Service]
ExecStartPre=/bin/sleep 1
TimeoutStartSec=15
